# sysconfig/systemd/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn/diematic.sock

[Install]
WantedBy=sockets.target
